# root/etc/s6-overlay/s6-rc.d/prowlarr/run
#!/command/with-contenv
/usr/bin/Prowlarr/Prowlarr -nobrowser -data=/config/prowlarr