# root/etc/s6-overlay/s6-rc.d/debug/run
#!/command/with-contenv

echo "=== DEBUG START ==="
echo "Current directory: $(pwd)"
echo "PATH: $PATH"
echo "Python version: $(python3 --version 2>/dev/null || echo 'Python not found')"
echo "UV version: $(uv --version 2>/dev/null || echo 'UV not found')"
echo ""

echo "Checking if Prowlarr binary exists:"
ls -la /usr/bin/Prowlarr/ || echo "Prowlarr directory not found"
echo ""

echo "Checking Comet directory:"
ls -la /app/comet/ || echo "Comet directory not found"
echo ""

echo "Checking uv sync status:"
cd /app/comet && uv run python -c "print('UV sync working')" 2>/dev/null || echo "UV sync failed"
echo ""

echo "=== DEBUG END ==="
