# root/etc/s6-overlay/s6-rc.d/cloudflared/run
#!/command/with-contenv
# Run cloudflared with the provided token
exec cloudflared tunnel --no-autoupdate run --token "$CF_TUNNEL_TOKEN"