# root/etc/s6-overlay/s6-rc.d/comet/run
#!/command/with-contenv

# Wait for Prowlarr to be ready
echo "Waiting for Prowlarr at http://127.0.0.1:9696..."
until curl -s http://127.0.0.1:9696 > /dev/null; do
  sleep 2
done
echo "Prowlarr is up. Starting Comet..."

# Start Comet
cd /app/comet
# FORCE unbuffered output so logs appear instantly
export PYTHONUNBUFFERED=1
# FORCE binding to 0.0.0.0 via command line args if possible, or env vars
export FASTAPI_HOST=0.0.0.0
export FASTAPI_PORT=8000

# Run uvicorn directly to bypass any main.py logic issues or gunicorn complexity
# This is the "Nuclear Option" to guarantee the server starts on the right port/interface
exec uvicorn comet.api.app:app --host 0.0.0.0 --port 8000 --workers 4